# liquidctl/liquiddyncfg.service
[Unit]
Description=Dynamic liquidctl adjustments

[Service]
Type=oneshot
ExecStart=/usr/local/bin/liquiddyncfg

[Install]
WantedBy=default.target
