# systemd/calvin-health-check-boot.timer
[Unit]
Description=Daily calvin hardware health summary
# Like the periodic timer, bound to the graphical session so it only fires while
# there is a session to read devices from and a bus to notify on. It reuses the
# boot service to re-emit the full "system state" summary once a day.
PartOf=graphical-session.target
After=graphical-session.target

[Timer]
OnCalendar=*-*-* 08:00:00

[Install]
WantedBy=graphical-session.target
