# kanata/launchd/com.jonasmalaco.karabiner-vhid-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jonasmalaco.karabiner-vhid-manager</string>

	<!-- Re-register (activate) the Karabiner DriverKit system extension at boot. This is
	     a one-shot: `activate` submits an OSSystemExtensionRequest and exits, so there is
	     no KeepAlive. Once the extension is user-approved (a manual GUI step, already done
	     on turing), macOS reloads it on its own each boot; running activate here is cheap
	     insurance that re-stages the current dext after a macOS or Karabiner update, which
	     is the usual cause of a "driver_version mismatched" failure. First-time approval
	     still needs the GUI and cannot be done headlessly. -->
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/.Karabiner-VirtualHIDDevice-Manager.app/Contents/MacOS/Karabiner-VirtualHIDDevice-Manager</string>
		<string>activate</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/var/log/karabiner-vhid-manager.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/karabiner-vhid-manager.log</string>
</dict>
</plist>
