# kanata/launchd/com.jonasmalaco.kanata.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jonasmalaco.kanata</string>

	<!-- kanata, reading the bombadil-linked config at ~/.config/kanata/kanata.kbd
	     (a symlink to dotfiles/kanata/kanata.kbd). Runs as root, like the manual
	     `sudo kanata --cfg ...` invocation it replaces. -->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/kanata</string>
		<string>--cfg</string>
		<string>/Users/jonas/.config/kanata/kanata.kbd</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<!-- Restart kanata ONLY on a non-zero exit (a crash); let a clean exit stay down.
	     This is launchd's equivalent of systemd's Restart=on-failure. Two payoffs:
	       - The kill switch (LCtrl+Space+Esc) exits 0 by default since kanata 1.11.0
	         (configurable via --emergency-exit-code; fixes #1902/#1904), so it now
	         really disables kanata and drops you to the default macOS layout for
	         diagnosis, instead of being relaunched within milliseconds.
	       - A sleep/wake crash (known open macOS bug #1357 on 1.11.0) exits non-zero
	         and IS relaunched, so the keyboard recovers on its own after waking.
	     ThrottleInterval caps relaunches to one per 10s, so a config that fails at
	     load retries calmly (visible in the log) rather than thrashing. kanata also
	     retries the Karabiner VHID connection itself, tolerating a cold start. -->
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>StandardOutPath</key>
	<string>/var/log/kanata.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/kanata.log</string>
</dict>
</plist>
