# kanata/launchd/local.kanata.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
    LaunchDaemon for kanata.

    Kanata needs root because the Karabiner-DriverKit-VirtualHIDDevice daemon exposes its IPC
    under `/Library/Application Support/org.pqrs/tmp/rootonly/`, which is mode 700 owned by
    root.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.kanata</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/kanata</string>
        <string>--cfg</string>
        <!--
            Root-owned copy (root:wheel 0644), not a user-writable file: kanata runs as
            root, so a config the login user could write would let any process running as
            that user escalate to root (the more so with `danger-enable-cmd`). Staged from
            the repo by `just kanata sync-config`.
        -->
        <string>/Library/Application Support/kanata/kanata.kbd</string>
        <!-- Don't wait for user input on exit, which would block automatic restart. -->
        <string>--no-wait</string>
        <!--
            Errors-only logging. At the default Info level kanata dumps the entire deflayer
            block on every layer switch, so /var/log/kanata.log grows without bound (and
            launchd has no log rotation; see "Logs" in ../CLAUDE.md for why naive rotation is
            worse than none). --quiet keeps config-load failures and other errors, which is
            what the log is actually for here, and drops the per-switch noise. To debug, swap
            for --debug (or remove this) and `just kanata reload local.kanata`.
        -->
        <string>--quiet</string>
        <!--
            Skip the 2s "release all keys" startup sleep. It guards against keys held during
            an interactive launch sticking; as a boot-time daemon there's no such launch, and
            ThrottleInterval already paces relaunches. kanata warns it "may cause keyboard
            issues on startup" in rare cases; drop this flag if a cold start ever misbehaves.
        -->
        <string>--nodelay</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!--
        Restart kanata ONLY on a non-zero exit (a crash); let a clean exit stay down.

        This is launchd's equivalent of systemd's Restart=on-failure. Two payoffs:

        - The kill switch (LCtrl+Space+Esc) exits 0 by default since kanata 1.11.0
        (configurable via --emergency-exit-code; fixes #1902/#1904), so it now really
        disables kanata and drops you to the default macOS layout for diagnosis, instead of
        being relaunched within milliseconds.

        - A sleep/wake crash (known open macOS bug #1357 on 1.11.0) exits non-zero and IS
        relaunched, so the keyboard recovers on its own after waking.

        ThrottleInterval caps relaunches to one per 10s, so a config that fails at load
        retries calmly (visible in the log) rather than thrashing. kanata also retries the
        Karabiner VHID connection itself, tolerating a cold start.
    -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/var/log/kanata.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/kanata.log</string>
</dict>
</plist>
